# supervisor/study/php-study.conf
[program:php-study-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/study/artisan queue:work
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/study/storage/logs/worker.log
stopwaitsecs=3600
